# airootfs/etc/profile.d/envs.sh
/home/shawn/dotfiles.git/main/extras/envs.sh